# etc/local.d/snail.start
# /etc/conf.d/local.start

snail.nv_pwr_off
